# bsp/meta-xiaopi/recipes-images/main/files/main.service
[Unit]
Description=Main program
After=initscript.service

[Service]
Type=simple
ExecStart=/usr/bin/main
Restart=always

[Install]
WantedBy=multi-user.target
