# bsp/meta-xiaopi/recipes-initscript/initscript/files/initscript.service
[Unit]
Description=start initscript upon first boot

[Service]
Type=simple
ExecStart=/bin/bash '/etc/rc.local'
